# root/etc/uci-defaults/90-luci-app-nginx-manager
#!/bin/sh
# SPDX-License-Identifier: AGPL-3.0-only

[ -n "${IPKG_INSTROOT}" ] && exit 0

mkdir -p /etc/nginx/conf.d/luci-manager/sites
mkdir -p /etc/nginx/conf.d/luci-manager/disabled
mkdir -p /etc/nginx/certs/luci-manager
mkdir -p /etc/nginx-manager/backups

cat > /etc/nginx/conf.d/luci-manager.conf <<'EOF'
# Generated by luci-app-nginx-manager.
# Source: /etc/config/nginx_manager
# SPDX-License-Identifier: AGPL-3.0-only

include /etc/nginx/conf.d/luci-manager/*.conf;
include /etc/nginx/conf.d/luci-manager/sites/*.conf;
EOF

uci -q get nginx_manager.global >/dev/null 2>&1 || {
	uci set nginx_manager.global=global
	uci set nginx_manager.global.managed_dir='/etc/nginx/conf.d/luci-manager'
	uci set nginx_manager.global.cert_dir='/etc/nginx/certs/luci-manager'
	uci set nginx_manager.global.backup_dir='/etc/nginx-manager/backups'
	uci set nginx_manager.global.auto_backup='1'
	uci set nginx_manager.global.test_before_reload='1'
	uci set nginx_manager.global.advanced_mode='0'
	uci set nginx_manager.global.dangerous_core_edit='0'
	uci set nginx_manager.global.max_backups='20'
	uci set nginx_manager.global.client_max_body_size=''
	uci set nginx_manager.global.keepalive_timeout=''
	uci set nginx_manager.global.gzip='0'
	uci set nginx_manager.global.server_tokens=''
	uci set nginx_manager.global.sendfile='0'
	uci set nginx_manager.global.access_log=''
	uci set nginx_manager.global.error_log=''
	uci set nginx_manager.global.ssl_protocols=''
	uci set nginx_manager.global.ssl_ciphers=''
	uci commit nginx_manager
}

chmod 600 /etc/config/nginx_manager 2>/dev/null

exit 0
